# .config/systemd/user/code.service
[Unit]
Description=Remote VSCode service

[Service]
ExecStart=code tunnel
WorkingDirectory=/home/z003tnse

[Install]
WantedBy=default.target
