# .config/sv/mpd/run
#!/bin/sh
mkdir -p /run/user/$(id -u)
exec /usr/bin/mpd --no-daemon ~/.mpd/mpd.conf
